# tftpd-dynamic/lib/systemd/system/tftpd-dynamic.socket
[Unit]
Description=Tftpd-dynamic Socket
PartOf=tftpd-dynamic.service

[Socket]
ListenDatagram=[::]:69

[Install]
WantedBy=sockets.target
